# browser_agency.ini
[program:browser_agency_1]
environment = LD_LIBRARY_PATH="/usr/local/lib:$LD_LIBRARY_PATH", PYTHONPATH="$PYTHONPATH:/usr/local/lib/python3.7/site-packages/browser_agency"
command = /usr/local/bin/python3.7 /usr/local/bin/run_browser_agency 
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /logs/browser_agency/browser_agency_1.log

[group:browser_agency]
programs=browser_agency_1
